# rpm-builder/SOURCES/ulmeter.timer
[Unit]
Description=Timer for ulmeter service

[Timer]
OnUnitActiveSec=60s
OnBootSec=10s
Unit=ulmeter.service

[Install]
WantedBy=multi-user.target
