# rpm-builder/SOURCES/ulmeter.service
[Unit]
Description = Get users load on host

[Service]
ExecStart=/usr/local/bin/ulmeter
