# src/etc/systemd/system/ulmtr-stat.service
[Unit]
Description = Stat total users load on host

[Service]
ExecStart=/usr/local/bin/ulmtr-stat
